# scaffolding
#!/bin/bash

### BEGIN INIT INFO
# Provides:          scaffolding
# Required-Start:    $syslog $network $named $time $local_fs
# Required-Stop:     $syslog $network $named $time $local_fs
# Default-Start:     5
# Default-Stop:      6
# Short-Description: scaffolding server for nodejs applications
# Description:       It is what it says one line above
### END INIT INFO

DIR=/home/pi/programs/nodejs/Scaffolding
DAEMON=server.coffee
DAEMON_NAME=scaffolding
PIDFILE=/var/run/$DAEMON_NAME.pid
LOGFILE=/var/log/$DAEMON_NAME.log

init() {
	rm $DIR/error.log
}

checkExecution() {
	if [ $? -eq 0 ]; then
		echo 'OK'
	else
		echo 'FAIL'
	fi
}

do_start() {
	echo "Starting $DAEMON_NAME service ..."
	echo -n "... Changing work directory ... "
	cd $DIR
	checkExecution

	echo -n "... Starting server ... "
	./$DAEMON > $LOGFILE &
	checkExecution

	echo -n "... Server started with PID: "
	echo $! > $PIDFILE
	echo "$!"
}

do_stop() {
	echo "Stopping $DAEMON_NAME service ..."
	echo -n "... killing process with PID: "
	APP_PID=$(<$PIDFILE)
	echo -n '$APP_PID '
	kill $APP_PID
	checkExecution
}

init
case $1 in
	start)
		do_start
		;;
	stop)
		do_stop
		;;
	status)
		;;
	restart)
		do_stop
		do_start
		;;
	install)
		echo "installing $DAEMON_NAME service ... "
		echo -n "... creating symlink ... "
		ln -s $DIR/$DAEMON_NAME /etc/init.d/$DAEMON_NAME 2>> $DIR/error.log
		checkExecution

		echo -n "... updating rc ... "
		update-rc.d $DAEMON_NAME defaults 2>> $DIR/error.log
		checkExecution
		;;
	uninstall|remove)
		echo "removing $DAEMON_NAME service ... "
		echo -n "... removing rc ... "
		update-rc.d $DAEMON_NAME remove 2>> $DIR/error.log
		checkExecution

		echo -n "... removing symlink ... "
		rm /etc/init.d/$DAEMON_NAME 2>> $DIR/error.log
		checkExecution
		;;
	*)
		echo 'Usage: service scaffolding start|stop|restart|status|install|remove|uninstall'
		echo 'install/remove requires admin rights'
		;;
esac

if [ -s "error.log" ]; then
	echo "Execution finished with errors, please check error.log file"
else
	echo "Execution finished successfully"
fi

exit 0
